# webhook.service
[Unit]
Description=GitHub Webhook Server for Auto-Deployment
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Resume-Analyzer
Environment="WEBHOOK_SECRET=your-secret-key-change-this"
ExecStart=/usr/bin/python3 /home/ubuntu/Resume-Analyzer/webhook_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
